# traffic-system.service
[Unit]
Description=Traffic Violations Management System
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/traffic-system
Environment="PORT=10000"
Environment="FLASK_ENV=production"
Environment="FLASK_DEBUG=false"
ExecStart=/usr/bin/gunicorn --config gunicorn_traffic_config.py traffic_app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
PrivateTmp=true
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
